# service/matchmaking.conf
description "Matchmaking server"
author "gordon.a.macpherson@gmail.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/matchmaking/BottleHoser/launch.sh